# bcachefs-wait@.service.in
[Unit]
Description=Every device in the Bcachefs filesystem with UUID %i is present

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=@sbindir@/bcachefs wait UUID=%i
